# general/.config/autostart/upwork.desktop
[Desktop Entry]
Comment=Upwork
Exec=/usr/bin/upwork
Hidden=true
Name=Upwork
StartupNotify=false
Terminal=false
Type=Application
Version=1.0
